# ops/launchd/templates/com.forme.ollama.autostart.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "https://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.forme.ollama.autostart</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/open</string>
    <string>-g</string>
    <string>-a</string>
    <string>/Applications/Ollama.app</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>LimitLoadToSessionType</key>
  <array>
    <string>Aqua</string>
  </array>

  <key>StandardOutPath</key>
  <string>__FORME_ROOT__/data/logs/ollama-autostart.log</string>

  <key>StandardErrorPath</key>
  <string>__FORME_ROOT__/data/logs/ollama-autostart.log</string>
</dict>
</plist>
